# systemd/system/enable-tv.service
[Unit]
Description=Enable HDMI output

[Service]
Type=oneshot
ExecStart=/bin/bash -c "echo detect > /sys/class/drm/card0-HDMI-A-1/status"
