# systemd/user/alsa-waybar.service
[Unit]
Description=ALSA Waybar Updater

[Service]
Type=simple
ExecStart=%h/dotfiles/bin/alsa-waybar

[Install]
WantedBy=default.target
